# .config/systemd/user/offlineimap.service
[Unit]

Description = fetch email with offlineimap

[Service]

Type = oneshot
ExecStart = notmuch new

[Install]

WantedBy = default.target
